# crontabs/pi4-root
@reboot /bin/bash /etc/iptables/add-openvpn-rules.sh